# cc-presentation.service
[Unit]
Description=CC Object Presentation
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cc-presentation
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=CORS_ORIGIN=http://72.56.96.129:3001
Environment=OPENAI_API_KEY=
ExecStart=/usr/bin/pnpm --filter server start
Restart=on-failure
RestartSec=10
MemoryMax=2G
CPUQuota=150%

[Install]
WantedBy=multi-user.target
